# roles/enecoq-data-fetcher/files/enecoq-data-fetcher.timer
[Unit]
Description=Enecoq data fetcher timer

[Timer]
OnCalendar=*:14,44
Persistent=true

[Install]
WantedBy=timers.target
