# roles/raspberrypi/files/backup-monthly-check.service
[Unit]
Description=Monthly restic repository integrity check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=RESTIC_REPOSITORY=/mnt/bak
Environment=RESTIC_PASSWORD_FILE=/etc/restic/password
Environment=RESTIC_PROGRESS_FPS=0.0166
Environment=XDG_CACHE_HOME=/var/cache/restic
ExecStart=/usr/local/sbin/backup-monthly-check
