# roles/raspberrypi/files/backup-daily.service
[Unit]
Description=Daily backup to NAS

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/backup-daily
