# roles/raspberrypi/files/backup-daily.timer
[Unit]
Description=Daily backup to NAS timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
