# roles/raspberrypi/files/backup-daily.service
[Unit]
Description=Daily backup to NAS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/backup-daily
