# roles/raspberrypi/files/network-failover.timer
[Unit]
Description=Network Failover Monitor Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
