# roles/raspberrypi/files/network-failover.service
[Unit]
Description=Network Failover Monitor
OnFailure=failure-notify@%n.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/network-failover
