# roles/raspberrypi/files/network-failover.service
[Unit]
Description=Network Failover Monitor

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/network-failover
