# roles/raspberrypi/files/network-failover.service
[Unit]
Description=Network Failover Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/network-failover
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
