# roles/ubuntu/files/failure-notify@.service
[Unit]
Description=Send failure notification for %i

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/failure-notify %i
