# roles/zabbix/mcp/files/zabbix-mcp-server.service
[Unit]
Description=Zabbix MCP Server
Documentation=https://github.com/initMAX/zabbix-mcp-server
After=network.target

[Service]
Type=simple
User=zabbix-mcp
Group=zabbix-mcp
ExecStart=/opt/zabbix-mcp/venv/bin/zabbix-mcp-server \
    --config /etc/zabbix-mcp/config.toml
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/zabbix-mcp/server.log
StandardError=append:/var/log/zabbix-mcp/server.log
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
ReadWritePaths=/var/log/zabbix-mcp

[Install]
WantedBy=multi-user.target
