# roles/enecoq/files/enecoq2mqtt.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
02 * * * * tats enecoq2mqtt
