# roles/mhz19/files/mhz192mqtt.service
[Unit]
Description=MH-Z19 CO2 sensor to MQTT

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mhz192mqtt
