# roles/nextcloud/files/nextcloud.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
*/5  *  *  *  * www-data php -f /srv/www/nextcloud/cron.php
