# roles/enecoq-data-fetcher/templates/enecoq-data-fetcher.cron.j2
# Data fetching job (twice daily with browser check)
14,44 * * * * root sleep $((RANDOM \% 60)) && (uvx enecoq-data-fetcher --email [% enecoq.email %] --password [% enecoq.password %] --period month --format json --output /srv/homeassistant/config/data/enecoq_data.json || (npx playwright install chromium --force && uvx enecoq-data-fetcher --email [% enecoq.email %] --password [% enecoq.password %] --period month --format json --output /srv/homeassistant/config/data/enecoq_data.json)) >> /var/log/enecoq-data-fetcher.log 2>&1
