# roles/sslcert/templates/certbot.cron.j2
#!/bin/bash

LOG_FILE="/var/log/certbot.log"

log() {
    echo "$(date '+%Y-%m-%d %H:%M:%S') - $1" >> "$LOG_FILE"
}

exec >> "$LOG_FILE" 2>&1

log "Starting certbot cron job"

source /opt/certbot/bin/activate

log "Updating certbot packages"
VIRTUAL_ENV=/opt/certbot uv pip install --upgrade certbot certbot-apache certbot-dns-route53

log "Renewing certificates"
certbot renew

if systemctl is-active --quiet apache2; then
    log "Restarting Apache"
    systemctl restart apache2
fi

log "Syncing certificates to Synology NAS"
/usr/local/sbin/sync-cert-to-synology

log "Certbot cron job completed"
