# roles/sslcert/templates/certbot.cron.j2
#!/bin/bash

. /opt/certbot/bin/activate

pip-review -a

sleep 10

certbot renew -q

systemctl restart apache2
