# roles/switchbotmeter/files/switchbotmeter2mqtt.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
01-59/5 * * * * root switchbotmeter2mqtt F5DAB9477E91 # Living Room
02-59/5 * * * * root switchbotmeter2mqtt F980AECF8644 # Bedroom
03-59/5 * * * * root switchbotmeter2mqtt FD13F5CA5563 # Storeroom
04-59/5 * * * * root switchbotmeter2mqtt EDE4EAA09EA3 # Balcony
