# roles/homeassistant-ha/secondary/files/homeassistant-failover.service
[Unit]
Description=Home Assistant failover health check
After=network-online.target docker.service
Wants=network-online.target
OnFailure=failure-notify@%n.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/homeassistant-failover
