# roles/homeassistant-ha/secondary/files/homeassistant-failover.timer
[Unit]
Description=Home Assistant failover health check (every 30 seconds)

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=5s
Unit=homeassistant-failover.service

[Install]
WantedBy=timers.target
