# roles/wordpress/files/wordpress.cron
#!/usr/bin/env bash

cd /srv/www/wordpress
./permset.sh