# roles/wordpress/files/wordpress.cron
#!/bin/sh

cd /srv/www/wordpress
./permset.sh
