# roles/wordpress/files/wordpress.cron
#!/bin/bash

# https://developer.wordpress.org/plugins/cron/hooking-wp-cron-into-the-system-task-scheduler/

wget -q --delete-after https://rewse.jp/blog/wp-cron.php
