# roles/postfix/files/certbot.cron
#!/bin/sh

certbot renew --force-renew --webroot-path=/srv/www/html --quiet

retval=$?

if [ $retval -ne 0 ]; then
  echo "<ERROR> certbot can't renew the certification"
  exit $retval
fi

systemctl reload httpd
systemctl reload postfix
systemctl reload dovecot
