# roles/sslcert/files/certbot.cron
#!/bin/bash

source /var/lib/python/certbot/bin/activate

pip-review -a

certbot renew -q

systemctl restart apache2
